# docker/ide/etc/cron/crontab
*/15 * * * * /opt/bin/backup >> /opt/log/backup.log 2>&1
